# recipes-local/epubd/epubd/epub-trigger.service
[Unit]
Description=Trigger upstream check for e-books to convert

[Service]
Type=oneshot
ExecStart=/usr/sbin/epub-trigger.sh
